# Linux/week3/DEBIAN/prerm
#!/bin/bash
#
# Prerm script - stops service before deleting

if [ -f "/etc/systemd/system/watchdog.service" ] ; then
  systemctl stop watchdog
  systemctl disable watchdog
fi
