# Linux/week3/etc/systemd/system/watchdog.service
[Unit]
Description=Standart notification service

[Service]
Type=notify
ExecStart=/usr/bin/watchdog_script

[Install]
WantedBy=multi-user.target
